# systemd/system/keepass_backup.service
[Unit]
Description=Backup keepass db to a different syncthing repo

[Service]
ExecStart=/usr/local/sbin/keepass_backup.sh

